# internal/control/install-edge.service
[Unit]
Description=simple_cdn edge agent
After=network-online.target nginx.service
Wants=network-online.target

[Service]
Type=simple
User=root
UMask=0027
EnvironmentFile=/opt/cdn-edge/config/edge.env
ExecStart=/opt/cdn-edge/bin/cdn-edge-agent
Restart=on-failure
RestartSec=5
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ReadWritePaths=/opt/cdn-edge
LimitNOFILE=65536

[Install]
WantedBy=multi-user.target
