# internal/control/install-edge-updater.service
[Unit]
Description=CDN Platform edge online upgrade %i
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
User=root
UMask=0077
EnvironmentFile=/opt/cdn-edge/config/edge.env
ExecStart=/opt/cdn-edge/bin/cdn-edge-agent upgrade-helper %i
TimeoutStartSec=35min
PrivateTmp=true
ProtectHome=true
NoNewPrivileges=true
